# setup/crontab
* * * * * php /www/var/tfmch_LARAVEL/artisan schedule:run > '/dev/null' 2>&1
0 0 * * * /www/var/tfmch_LARAVEL/update.sh
